# dpkg/elementalclinic/debian/cron.d
#
# Regular cron jobs for the elementalclinic package
#
0  *  *  *  *  /usr/share/elementalclinic/bin/send_notifications.pl
